# stage-tdc/02-systemd-units/files/vena-pulse-onboarding.service
[Unit]
Description=Vena Pulse first-boot captive portal (run only if no config)
Documentation=https://github.com/Chaotic-Lab/caos-lab-sbc
ConditionPathExists=!/etc/vena-pulse/config.json
After=network.target
Conflicts=vena-pulse-scanner-wifi.service vena-pulse-scanner-ble.service vena-pulse-publisher.service

[Service]
Type=simple
User=root
ExecStart=/usr/local/bin/vena-pulse-onboarding
Restart=on-failure
RestartSec=10s
TimeoutStopSec=15s

[Install]
WantedBy=multi-user.target
