# stage-tdc/02-systemd-units/files/vena-pulse-update.service
[Unit]
Description=Vena Pulse — OTA self-update (Sprint E deliverable; stub on this image)
# chaoticlab-sbc is PRIVATE, so an anonymous `pip install git+https://...` would
# fail on a deployed Pi. Real OTA is a Sub-plan E deliverable: pull the
# cosign-verified wheel from MinIO (s3.chaoticlab.cloud/datapool-tdc/firmware/
# pulse/...), verify-blob against the chaoticlab-sbc workflow identity, then
# install. See spec §20 R10 + master plan cross-cutting concerns. Until then
# this is a logged no-op (and is gated below so it never runs pre-onboarding).
ConditionPathExists=/etc/vena-pulse/config.json
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/bin/echo "Vena Pulse OTA self-update is a Sub-plan E deliverable (cosign-verified wheel from MinIO). No-op on this image."
TimeoutStartSec=60s

[Install]
WantedBy=multi-user.target
